# deploy/pik-scan-dev.service
[Unit]
Description=Daily price scan for non-PIK developers (ФСК, Донстрой, А101, Level, Абсолют, MR Group, Гранель, Инград, Брусника)
After=network.target pik-scan.service
# Defensive: если кто-то запустит оба `systemctl start pik-scan.service
# pik-scan-dev.service` одновременно, After= даст systemd ordering hint.
# Основной триггер — OnSuccess/OnFailure в pik-scan.service.

[Service]
Type=oneshot
User=pik
Group=pik
WorkingDirectory=/opt/pik
Environment=TZ=Europe/Moscow
# Запускается отдельно от bin.scan (ПИК): другой набор источников,
# та же БД. --all обходит всех застройщиков из реестра параллельно.
ExecStart=/opt/pik/venv/bin/python -m bin.scan_dev --db /opt/pik/data/pik.db --all
StandardOutput=journal
StandardError=journal
