# deploy/pik-scan.service
[Unit]
Description=Daily price scan — all 10 developers (PIK + 9 others) → SQLite
After=network.target
# До R5 рефактора (2026-05-25) было два юнита: pik-scan (только ПИК) и
# pik-scan-dev (только не-ПИК), скованных OnSuccess-цепью. После R1+R5
# scan_dev обрабатывает все 10 застройщиков одной командой — двух юнитов
# больше не нужно. install.sh удаляет pik-scan-dev.{service,timer} на
# upgrade. Имя «pik-scan» оставлено как label — деплоймент-имя.

[Service]
Type=oneshot
User=pik
Group=pik
WorkingDirectory=/opt/pik
Environment=TZ=Europe/Moscow
# --all обходит всех 10 застройщиков из реестра pik.sources.* параллельно
# (--workers=6). ПИК через bin.sources.pik (новый адаптер) — никаких
# отдельных команд.
ExecStart=/opt/pik/venv/bin/python -m bin.scan_dev --db /opt/pik/data/pik.db --all
StandardOutput=journal
StandardError=journal

# Resource limits. Scan держит до 100МБ под буферы парсинга (68 ЖК
# параллельно), 512МБ — комфортный запас x5. CPUQuota 400% = до 4-х
# воркеров на полной отдаче (ThreadPoolExecutor max_workers=6).
MemoryMax=512M
CPUQuota=400%
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
NoNewPrivileges=yes
ReadWritePaths=/opt/pik/data
TimeoutStartSec=1800
# 30мин cap: при здоровых API скан занимает 2-5мин; жёсткая стенка на
# случай зависшего HTTP-соединения, которое requests-timeout не поймал.
