# deploy/pik-scan-dev.timer
[Unit]
Description=Run pik-scan-dev daily at 06:30 MSK
Requires=pik-scan-dev.service

[Timer]
# 06:30 — на полчаса позже ПИК-скана (pik-scan.timer в 06:00), чтобы
# два писателя не делили запись в одну БД на пике.
OnCalendar=*-*-* 06:30:00 Europe/Moscow
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
