# deploy/pik-scan.timer
[Unit]
Description=Run pik-scan daily at 06:00 MSK
Requires=pik-scan.service

[Timer]
OnCalendar=*-*-* 06:00:00 Europe/Moscow
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
