# deploy/pik.service
[Unit]
Description=Datasette serving PIK Narvin snapshots
After=network.target

[Service]
Type=simple
User=pik
Group=pik
WorkingDirectory=/opt/pik
Environment=TZ=Europe/Moscow
ExecStart=/opt/pik/venv/bin/datasette serve /opt/pik/data/pik.db \
  --metadata /opt/pik/metadata.yml \
  --host 127.0.0.1 --port 5051 \
  --setting truncate_cells_html 0 \
  --setting max_returned_rows 5000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
